# c920.rules
ATTRS{idVendor}=="046d", ATTRS{idProduct}=="0892", RUN+="/usr/local/bin/c920_focus.sh"

